# linux/.config/systemd/refresh-gh-cache.timer
[Unit]
Description=Run my script every 15 minutes
Requires=myscript.service

[Timer]
Unit=myscript.service
# Run every 15 minutes
OnUnitActiveSec=15min
# Also run 15 minutes after boot
OnBootSec=15min
# Random delay to distribute system load (optional)
RandomizedDelaySec=30sec

[Install]
WantedBy=timers.target
