# linux/.config/systemd/refresh-gh-cache.service
[Unit]
Description=Run my script every 15 minutes
After=network.target

[Service]
Type=oneshot
ExecStart=/path/to/your/script.sh
# If your script requires specific environment variables:
# Environment=VAR1=value1 VAR2=value2

[Install]
WantedBy=default.target
